# systemd/netconsole2journal.socket
[Unit] 
Description=Netconsole UDP receiver

[Socket] 
ListenDatagram=2

[Install] 
WantedBy=sockets.target
